# chef-bookbooks/O/oozie/templates/default/oozie.service.erb
[Unit]
Description = Oozie Server
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.oozie.user %>

PIDFile = <%= node.oozie.home %>/oozie-server/temp/oozie.pid
ExecStart = <%= node.oozie.home %>/bin/oozie-start.sh
ExecStop = <%= node.oozie.home %>/bin/oozie-stop.sh
LimitNOFILE = 10000

Type = forking

[Install]
WantedBy = multi-user.target
