# chef-bookbooks/H/hypercontainer/templates/default/etc/systemd/system/hyperd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/hyperd <%= @hyperd_opts.join(' ') %>
