# chef-bookbooks/S/spark-platform/templates/default/systemd/spark-master.service.erb
[Unit]
Description=Spark Master
After=network.target

[Service]
Type=forking
User=<%= @user %>
Group=<%= @group %>
ExecStart="<%= @prefix_root %>/spark/sbin/start-master.sh"
Restart=on-failure

[Install]
WantedBy=multi-user.target
