# chef-bookbooks/S/spark-platform/templates/default/systemd/spark-worker.service.erb
[Unit]
Description=Spark Slave
After=network.target

[Service]
Type=forking
User=<%= @user %>
Group=<%= @group %>
ExecStart=<%= @prefix_root %>/spark/sbin/start-slave.sh <%= @master_url %>
Restart=on-failure

[Install]
WantedBy=multi-user.target
