# chef-bookbooks/P/pg/templates/default/postgresql.service.erb
.include /usr/lib/systemd/system/postgresql.service
[Service]
Environment=
Environment=PGPORT=<%= node['pg']['config']['server']['port'] %>
Environment=PGDATA=<%= node['pg']['datadir'] %>
