# chef-bookbooks/D/docker-grid/templates/default/etc/systemd/system/docker.service.d/override.conf
<%
dockerd = '/usr/bin/dockerd'
# old daemon command
dockerd = '/usr/bin/docker daemon' unless File.exist?(dockerd)
-%>
[Service]
ExecStart=
ExecStart=<%= dockerd %> <%= @docker_opts.join(' ') %>
